# exercises/09-services-templates-and-registering-variables/app.service.j2
[Unit]
Description=Application
After=network.target

[Service]
Type=simple
Restart=always
User={{app_user}}
Group={{app_group}}
Environment=PORT={{app_port}}
ExecStart={{node_path.stdout}} {{app_root}}/app.js
